# cmd/tailscaled/tailscale-online.target
[Unit]
Description=Tailscale is online
Requires=tailscale-wait-online.service
After=tailscale-wait-online.service
